# site-modules/abide_linux/templates/systemd/aidecheck.timer.epp
<%- | String $timer_schedule | -%>
[Unit]
Description=Aide check every day at midnight

[Timer]
OnCalendar=<%= $timer_schedule %>
Unit=aidecheck.service

[Install]
WantedBy=multi-user.target
